# templates/circupon-boot.conf.erb
start on runlevel [2345]

script
  IP=$(ifconfig | grep "inet addr:" | grep -v "127.0.0.1" | awk '{print $2}' | cut -f2 -d:)
  CHANGE_BATCH=$(cat <<@@
{ "Changes": [ { "Action": "UPSERT", "ResourceRecordSet": { "Name": "<%= fqdn %>.", "Type": "A", "TTL": 300, "ResourceRecords": [ { "Value": "$IP" } ] } } ] }
@@
)

  hostname <%= fqdn %>
  aws route53 change-resource-record-sets --hosted-zone-id ZNEX827OEXWWQ --change-batch "$CHANGE_BATCH"
end script

