# templates/circupon-app-update.conf.erb
start on circupon-update

normal exit 0 2

script
  puppet agent --test
end script
