# templates/circupon-play-app.conf.erb
start on (stopped circupon-play-<%= appName %>-deploy or runlevel [2345])

stop on runlevel [!2345]

respawn

respawn limit 10 5

script
  cd /usr/local/circupon/<%= appName %>
  if [ -x current/bin/<%= appName %> ] 
  then 
    export JAVA_OPTS="-javaagent:/usr/local/circupon/<%= appName %>/newrelic/newrelic.jar -DdeploymentEnv=<%= @deploymentEnv %>"
    ./current/bin/<%= appName %> 
  fi
end script
