# files/circupon-update.conf
start on circupon-update

export FACTER_deployment_env

script
  puppet agent --test --server puppet.internal.pirc.com
end script
