# files/circupon-update.conf
start on circupon-update

script
  puppet agent --test --server puppet.internal.pirc.com
end script
