# files/circupon-update.conf
start on circupon-update

script
  puppet agent --test
end script
